# launchd/com.xwmx.nb-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily nb sync on macOS hosts (the notebook's ~/.nb/home/CLAUDE.md Sync
  section refers to this). On weddle the pkb-agent timer syncs instead.

  Install:
    cp ~/.dotfiles/launchd/com.xwmx.nb-sync.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.xwmx.nb-sync.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xwmx.nb-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>nb sync --all</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/ben/Library/Logs/nb-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ben/Library/Logs/nb-sync.log</string>
</dict>
</plist>
